# infra/vast/supervisor/conf.d/zkpow-prover.conf
[program:zkpow-prover]
environment=PROC_NAME="%(program_name)s"
command=/opt/supervisor-scripts/zkpow-prover.sh
autostart=true
; Don't auto-restart — if a batch fails investigate before restarting manually:
;   supervisorctl restart zkpow-prover
autorestart=false
startretries=3
; Must stay alive for 60s to be considered "started" (avoids depleting startretries
; during the wait-for-db loop which exits quickly if DB is already present)
startsecs=60
stopasgroup=true
killasgroup=true
stopsignal=TERM
; 5 min grace so an in-flight batch has time to reach a checkpoint before SIGKILL
stopwaitsecs=300
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_events_enabled=true
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
